# openSeSSHIAMe/etc/openSeSSHIAMe-oneshot.timer
[Unit]
Description=Run openSeSSHIAMe-oneshot periodically

[Timer]
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
